# .docker/development/files/etc/cron.d/plugins.glpi-project.org
*/5 * * * *    /usr/local/bin/php /var/www/plugins.glpi-project.org/misc/run_tasks.php &>/dev/null
